# node_relay/ops/node-relay-poller.service
[Unit]
Description=node-relay poller (remote node side) — relay pending/ -> execute -> terminal/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust User/WorkingDirectory/venv path per node. No vendor-specific defaults.
User=relay
WorkingDirectory=/home/relay/node-relay
EnvironmentFile=/home/relay/.config/node-relay/relay.env
ExecStart=/home/relay/node-relay/.venv/bin/python -m node_relay.node_poller --interval 10 --executor local+cloud
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
